# config/s6-overlay/s6-rc.d/create-buckets/up
/etc/s6-overlay/scripts/create-buckets.sh
